# microclaw.service
[Unit]
Description=Microclaw AI Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/microclaw
ExecStart=/home/ubuntu/microclaw/target/release/microclaw start
Restart=always
RestartSec=3

# 日志
StandardOutput=append:/var/log/microclaw.log
StandardError=append:/var/log/microclaw.log

# 环境变量
Environment=RUST_LOG=info
Environment=HOME=/home/ubuntu

# 优雅关闭等待
TimeoutStopSec=10
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
