# firecracker-saas/systemd/microclaw-nginx.service
[Unit]
Description=MicroClaw SaaS Nginx Reverse Proxy
After=network-online.target microclaw-control.service
Wants=network-online.target

[Service]
Type=forking
PIDFile=/run/nginx.pid
ExecStartPre=/usr/sbin/nginx -t -c /etc/nginx/conf.d/microclaw-saas.conf
ExecStart=/usr/sbin/nginx -c /etc/nginx/conf.d/microclaw-saas.conf
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s QUIT $MAINPID
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
