# firecracker-saas/guest/microclaw.service
[Unit]
Description=MicroClaw AI Chat Bot
After=network.target

[Service]
Type=simple
User=microclaw
Group=microclaw
WorkingDirectory=/app
Environment=MICROCLAW_CONFIG=/data/config/config.yaml
Environment=RUST_LOG=info
ExecStart=/usr/local/bin/microclaw start --config /data/config/config.yaml
Restart=on-failure
RestartSec=2
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
